# os-core/systemd/junction-reconcile-role.path
[Unit]
Description=Watch cluster.json and reconcile Docker role immediately

[Path]
PathChanged=/etc/junction/cluster.json

[Install]
WantedBy=multi-user.target
